# packaging/deb/debian/postinst
#!/bin/sh
# CloudWorkstation post-installation script
# This script is run after the package is installed

set -e

# Source debconf library for configuration management
. /usr/share/debconf/confmodule

# Package name
PACKAGE="cloudworkstation"
USER="cloudworkstation"
GROUP="cloudworkstation"
HOME_DIR="/var/lib/cloudworkstation"
LOG_DIR="/var/log/cloudworkstation"
CONFIG_DIR="/etc/cloudworkstation"

case "$1" in
    configure)
        # Create cloudworkstation system user and group if they don't exist
        if ! getent group "$GROUP" >/dev/null 2>&1; then
            addgroup --system "$GROUP"
        fi
        
        if ! getent passwd "$USER" >/dev/null 2>&1; then
            adduser --system --ingroup "$GROUP" --home "$HOME_DIR" \
                    --no-create-home --disabled-login \
                    --shell /usr/sbin/nologin \
                    --gecos "CloudWorkstation System User" "$USER"
        fi
        
        # Create and set ownership of directories
        mkdir -p "$HOME_DIR" "$LOG_DIR" "$CONFIG_DIR/aws"
        
        # Set proper ownership
        chown "$USER:$GROUP" "$HOME_DIR"
        chown "$USER:$GROUP" "$LOG_DIR"
        chown root:"$GROUP" "$CONFIG_DIR"
        chown root:"$GROUP" "$CONFIG_DIR/aws"
        
        # Set proper permissions
        chmod 755 "$HOME_DIR"
        chmod 750 "$LOG_DIR"
        chmod 750 "$CONFIG_DIR"
        chmod 750 "$CONFIG_DIR/aws"
        
        # Set permissions on config files
        if [ -f "$CONFIG_DIR/daemon.conf" ]; then
            chown root:"$GROUP" "$CONFIG_DIR/daemon.conf"
            chmod 640 "$CONFIG_DIR/daemon.conf"
        fi
        
        if [ -f "$CONFIG_DIR/aws/config.template" ]; then
            chown root:"$GROUP" "$CONFIG_DIR/aws/config.template"
            chmod 640 "$CONFIG_DIR/aws/config.template"
        fi
        
        if [ -f "$CONFIG_DIR/aws/credentials.template" ]; then
            chown root:"$GROUP" "$CONFIG_DIR/aws/credentials.template"
            chmod 640 "$CONFIG_DIR/aws/credentials.template"
        fi
        
        # Enable systemd service
        if [ -d /run/systemd/system ]; then
            systemctl --system daemon-reload >/dev/null || true
            if ! systemctl is-enabled cloudworkstation.service >/dev/null 2>&1; then
                systemctl enable cloudworkstation.service >/dev/null || true
            fi
            
            # Create systemd drop-in directory for local overrides
            mkdir -p /etc/systemd/system/cloudworkstation.service.d
        fi
        
        # Display installation success message
        cat << 'EOF'

â•”â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•—
â•‘                   CloudWorkstation Successfully Installed                    â•‘
â• â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•£
â•‘                                                                              â•‘
â•‘  ðŸŽ‰ CloudWorkstation v0.4.2 has been installed successfully!                â•‘
â•‘                                                                              â•‘
â•‘  ðŸ“‹ What's installed:                                                        â•‘
â•‘     âœ… Core binaries: cws, cwsd                                              â•‘
â•‘     âœ… Systemd service: cloudworkstation.service                             â•‘
â•‘     âœ… Configuration: /etc/cloudworkstation/                                 â•‘
â•‘     âœ… User account: cloudworkstation (system user)                          â•‘
â•‘                                                                              â•‘
â•‘  ðŸš€ Next steps:                                                              â•‘
â•‘                                                                              â•‘
â•‘     1. Configure AWS credentials:                                            â•‘
â•‘        sudo cp /etc/cloudworkstation/aws/credentials.template \             â•‘
â•‘               /etc/cloudworkstation/aws/credentials                          â•‘
â•‘        sudo cp /etc/cloudworkstation/aws/config.template \                  â•‘
â•‘               /etc/cloudworkstation/aws/config                               â•‘
â•‘        sudo nano /etc/cloudworkstation/aws/credentials                      â•‘
â•‘                                                                              â•‘
â•‘     2. Start the service:                                                    â•‘
â•‘        sudo systemctl start cloudworkstation                                â•‘
â•‘                                                                              â•‘
â•‘     3. Test the installation:                                                â•‘
â•‘        cws --version                                                         â•‘
â•‘        cws templates                                                         â•‘
â•‘                                                                              â•‘
â•‘  ðŸ“š Documentation: https://github.com/scttfrdmn/prism            â•‘
â•‘  ðŸ› Issues: https://github.com/scttfrdmn/prism/issues            â•‘
â•‘                                                                              â•‘
â•‘  ðŸ’¡ Quick start guide:                                                       â•‘
â•‘     Run 'man cws' for detailed usage instructions                           â•‘
â•‘     Run 'man cwsd' for daemon configuration help                            â•‘
â•‘                                                                              â•‘
â•šâ•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•

EOF
        ;;
        
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
        
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# Automatically added by dh_systemd_enable/11.1.6ubuntu2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
    # This will only remove masks created by d-s-h on package removal.
    deb-systemd-helper unmask 'cloudworkstation.service' >/dev/null || true

    # was-enabled defaults to true, so new installations run enable.
    if deb-systemd-helper --quiet was-enabled 'cloudworkstation.service'; then
        # Enables the unit on first installation, creates new
        # symlinks on upgrades if the unit file has changed.
        deb-systemd-helper enable 'cloudworkstation.service' >/dev/null || true
    else
        # Update the statefile to add new symlinks (if any), which need to be
        # cleaned up on purge. Also remove old symlinks.
        deb-systemd-helper update-state 'cloudworkstation.service' >/dev/null || true
    fi
fi
# End automatically added section

exit 0